# scripts/file-tree.timer
[Unit]
Description=Run file-tree once every 30 min and on boot

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
